# vault/rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Vault ingress proxy – finish
# ==============================================================================
bashio::log.warning "nginx exited"
